# homedir/crontab
* * * * * cd /home/user/GDrive/googleAccount;grive -s pathToSync    # or without "-s pathToSync" to sync all
