# homedir/crontab
* * * * * cd /home/user/GDrive/googleAccount;grive -s pathToSync
